# nuage.conf
[program:nuage]
command = gunicorn -b 0.0.0.0:8000 --workers=3 script:app
directory = /home/user/app
