# nuage.conf
[program:nuage]
#command = newrelic-admin run-program gunicorn  -b 0.0.0.0:5555 wsgi:app
command = python script.py
directory = /app
